# initscripts/maat.service.in
[Unit]
Description=Maat Attestation Manager
After=syslog.target network.target

[Service]
TimeoutStartSec=30
EnvironmentFile=/etc/default/maat
ExecStart=@bindir@/attestmgr -C @sysconfdir@/attestmgr-config.xml
KillMode=process
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
